# deploy/traefik/systemd/scu-laderaum-planer-update.service
[Unit]
Description=Update SCU-Laderaum-Planer from latest GitHub release
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/scu-laderaum-planer-update
ExecStart=/bin/bash -lc 'exec "${DEPLOY_DIR:?Please set DEPLOY_DIR in /etc/default/scu-laderaum-planer-update}/deploy-latest.sh"'
